# install/templates/halix-api.service
[Unit]
Description=Halix Cloud API (NestJS)
After=network.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=halix
Group=halix
WorkingDirectory=__INSTALL_ROOT__/apps/api
EnvironmentFile=__INSTALL_ROOT__/apps/api/.env
ExecStart=/usr/bin/env node dist/main.js
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
